# root/etc/services.d/queue-cleanarr/run
#!/usr/bin/with-contenv bash

s6-setuidgid ${PUID}:${PGID} python3 /usr/bin/queue-cleanarr.py --config "/config/config.xml" --interval 1800
